# config/supervisor/fcserver.conf
[program:fcserver]
command=/usr/local/bin/fcserver /etc/fcserver.json
autostart=true
autorestart=true
startretries=10
logfile=/var/log/fcserver.log
