# config/supervisor/structure_input.conf
[program:structure_input]
command=/home/pi/repo/Dreamland/structure_input.py /etc/dreamland.py
autostart=true
autorestart=true
startretries=10
