# config/supervisor/structure_output.conf
[program:structure_output]
command=/home/pi/repo/Dreamland/structure_output.py /etc/dreamland.py
autostart=true
autorestart=true
startretries=10
